# deploy/alpaca-scheduler.service
[Unit]
Description=Alpaca DUOL paper trading scheduler
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/alpaca-realtime-trading
Environment=PATH=/path/to/alpaca-realtime-trading/venv/bin
ExecStart=/path/to/alpaca-realtime-trading/venv/bin/python -u scheduler.py
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
